# infra/systemd/praxis-scheduler.service
[Unit]
Description=praxis-v2 scheduler + dead-man switch
After=network.target postgresql.service
Wants=postgresql.service
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=180
User=avyuk
Group=avyuk
WorkingDirectory=/home/avyuk/dev/praxis-v2
EnvironmentFile=/home/avyuk/dev/praxis-v2/.env
ExecStart=/home/avyuk/dev/praxis-v2/.venv/bin/python -m services.scheduler.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
