# infra/systemd/praxis-dashboard.service
[Unit]
Description=praxis-v2 dashboard
After=network.target postgresql.service

[Service]
Type=simple
User=praxis
Group=praxis
WorkingDirectory=/opt/praxis-v2
EnvironmentFile=/etc/praxis/praxis.env
ExecStart=/opt/praxis-v2/.venv/bin/python -m services.dashboard.app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
